# metin2pricechart.service.example
# Przykładowa jednostka systemd – automatyczny start i restart aplikacji.
# Instalacja na VPS:
#   1. W katalogu projektu: python3 -m venv venv && ./venv/bin/pip install -r requirements.txt
#   2. Skopiuj: sudo cp metin2pricechart.service.example /etc/systemd/system/metin2pricechart.service
#   3. Edytuj: User, WorkingDirectory, ExecStart (użyj venv jeśli masz venv!)
#   4. Włącz i start: sudo systemctl daemon-reload && sudo systemctl enable --now metin2pricechart

[Unit]
Description=Metin2 Price Chart
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Metin2PriceChart
Environment=PORT=5001
Environment=DATABASE_PATH=/root/Metin2PriceChart/price_history.db
# Secret do webhooka i /api/logs – dowolny długi losowy string (bez spacji lub w cudzysłowach)
# Environment=DEPLOY_WEBHOOK_SECRET=twoj_tajny_losowy_string
# Environment=DEPLOY_SERVICE=metin2pricechart
# Opcjonalnie: LOW_MEMORY=1, LOG_FILE=/var/log/metin2pricechart/app.log
# WAŻNE: Użyj Pythona z venv (gdzie jest flask) – inaczej: No module named 'flask'
ExecStart=/root/Metin2PriceChart/venv/bin/python main.py
# Gdy NIE używasz venv (pakiety globalnie: pip3 install -r requirements.txt):
# ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
